# system/etc/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

runuser -l pi -c '/home/pi/picopterx/scripts/mjpg-streamer.sh'
runuser -l pi -c 'screen -dmS test bash -c "sudo /home/pi/picopterx/scripts/run-server.sh; exec bash"'
exit 0
